# openwrt/etc/init.d/set_vpn_rules
#!/bin/sh /etc/rc.common
# (C) 2022 nearhop.com

# Set VPN rules
START=12
boot() {
  /sbin/create_vpn_forward_firewall.sh
}

